# systemd/hillei_hw.slice
[Unit]
Description=Hillel_HW Resource Slice
Before=slices.target

[Slice]
CPUQuota=25%

MemoryMax=100M
